# scripts/limen-agent.service
[Unit]
Description=LIMEN Agent Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=darc0
Group=darc0
WorkingDirectory=/home/darc0/projects/LIMEN/backend/agent
ExecStart=/home/darc0/projects/LIMEN/backend/agent/target/release/agent
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=limen-agent
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

