# scripts/limen-backend.service
[Unit]
Description=LIMEN Backend Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=darc0
Group=darc0
WorkingDirectory=/home/darc0/projects/LIMEN/backend
ExecStart=/home/darc0/projects/LIMEN/backend/server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=limen-backend
EnvironmentFile=/home/darc0/projects/LIMEN/backend/.env
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target




